# services/haproxy.service.conf
[program:haproxy]
user=root
command=haproxy -f /storage/haproxy.cfg
stdout_events_enabled = true
stderr_events_enabled = true
